# company-backend.service
[Unit]
Description=Company Backend Spring Boot Application
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/company-backend
ExecStart=/usr/bin/java -jar /opt/company-backend/company_backend-0.0.1-SNAPSHOT.jar --spring.profiles.active=prod-aws
Restart=always
RestartSec=10

# Environment variables
Environment=JAVA_HOME=/usr/lib/jvm/java-17-amazon-corretto
Environment=SSL_KEYSTORE_PASSWORD=changeit

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=company-backend

[Install]
WantedBy=multi-user.target